# ansible/roles/haproxy/files/etc/systemd/system/haproxy.service.d/slice.conf
# Assign HAProxy to the haproxy slice.
[Service]
Slice=haproxy.slice
